# etc/uci-defaults/luci-i18n-ddns-sk
uci set luci.languages.sk='SlovenÄina (Slovene)'; uci commit luci
