# etc/uci-defaults/luci-i18n-qos-ca
uci set luci.languages.ca='Català (Catalan)'; uci commit luci
