# etc/uci-defaults/luci-i18n-openvpn-sv
uci set luci.languages.sv='Svenska (Swedish)'; uci commit luci
