# etc/uci-defaults/luci-i18n-tinyproxy-pt
uci set luci.languages.pt='Português (Portuguese)'; uci commit luci
