# etc/uci-defaults/luci-i18n-multiwan-cs
uci set luci.languages.cs='Čeština (Czech)'; uci commit luci
