# etc/uci-defaults/luci-i18n-openvpn-zh-tw
uci set luci.languages.zh_tw='臺灣華語 (Taiwanese)'; uci commit luci
