# etc/uci-defaults/luci-i18n-openvpn-pt-br
uci set luci.languages.pt_br='Português do Brasil (Brazialian Portuguese)'; uci commit luci
