# etc/uci-defaults/luci-i18n-openvpn-vi
uci set luci.languages.vi='Tiếng Việt (Vietnamese)'; uci commit luci
