# etc/uci-defaults/luci-i18n-qos-pl
uci set luci.languages.pl='Polski (Polish)'; uci commit luci
