# etc/uci-defaults/luci-i18n-openvpn-uk
uci set luci.languages.uk='украї́нська (Ukrainian)'; uci commit luci
