# etc/uci-defaults/luci-i18n-qos-tr
uci set luci.languages.tr='Türkçe (Turkish)'; uci commit luci
