# etc/uci-defaults/luci-i18n-qos-ro
uci set luci.languages.ro='Română (Romanian)'; uci commit luci
