# etc/uci-defaults/luci-i18n-ddns-no
uci set luci.languages.no='Norsk (Norwegian)'; uci commit luci
