# etc/uci-defaults/luci-i18n-tinyproxy-en
uci set luci.languages.en='English'; uci commit luci
