# etc/uci-defaults/luci-i18n-multiwan-el
uci set luci.languages.el='Ελληνικά (Greek)'; uci commit luci
