# etc/uci-defaults/luci-i18n-base-ja
uci set luci.languages.ja='日本語 (Japanese)'; uci commit luci
