# etc/uci-defaults/luci-i18n-qos-ms
uci set luci.languages.ms='Bahasa Melayu (Malay)'; uci commit luci
