# etc/uci-defaults/luci-theme-itusnetworks
#!/bin/sh
uci batch <<-EOF
	set luci.themes.ItusNetworks=/luci-static/itusnetworks
	commit luci
EOF
exit 0
