# etc/uci-defaults/luci-i18n-multiwan-ru
uci set luci.languages.ru='Русский (Russian)'; uci commit luci
