# etc/uci-defaults/luci-i18n-qos-de
uci set luci.languages.de='Deutsch (German)'; uci commit luci
