# etc/uci-defaults/luci-i18n-qos-fr
uci set luci.languages.fr='Français (French)'; uci commit luci
