# etc/uci-defaults/luci-i18n-firewall-hu
uci set luci.languages.hu='Magyar (Hungarian)'; uci commit luci
