# etc/uci-defaults/luci-i18n-base-es
uci set luci.languages.es='Español (Spanish)'; uci commit luci
