# etc/uci-defaults/luci-i18n-qos-it
uci set luci.languages.it='Italiano (Italian)'; uci commit luci
