# usr/lib/opkg/info/base-files.list
/etc/e2guardian/lists/blacklists/drugs/urls
/etc/e2guardian/lists/blacklists/sexuality/domains
/etc/rc.local_RESTORE
/etc/openwrt_version
/etc/config/network-router
/etc/banner.failsafe
/etc/resolv.conf
/etc/init.d/done
/etc/csamd.conf
/bin/config_generate
/var
/etc/itus-setup.sh_ROUTER
/lib/preinit/00_preinit.conf
/etc/init.d/system
/etc/preinit
/etc/e2guardian/lists/blacklists/racism/urls
/www/demo/index.html
/www/demo/usage.htm
/usr/lib/libaivsec.so.5
/etc/shells
/lib64
/lib/preinit/70_initramfs_test
/etc/e2guardian/lists/blacklists/racism/expressions
/etc/device_info
/usr/lib64
/etc/protocols
/etc/uci-defaults/12_network-generate-ula
/usr/lib/libaivse000.so
/etc/e2guardian/lists/blacklists/ads/expressions
/lib/preinit/10_indicate_preinit
/etc/upgrade.sh
/lib/functions/uci-defaults-new.sh
/etc/e2guardian/lists/blacklists/piracy/domains
/etc/itus-setup.sh
/www/demo/itus.jpg
/etc/config/network-slave
/etc/e2guardian/lists/blacklists/phishing/urls
/etc/config/system
/etc/e2guardian/lists/blacklists/malware/domains
/etc/itus-setup.sh_GATEWAY
/usr/lib/libaivsec.so
/etc/config/network-master
/etc/fstab
/etc/sysctl.conf
/sbin/sysupgrade
/lib/preinit/99_10_run_init
/etc/config/firewall-router
/sbin/led.sh
/etc/e2guardian/lists/blacklists/proxies/domains
/usr/share/file/magic
/lib/preinit/50_indicate_regular_preinit
/etc/banner
/etc/suricata/threshold.config
/etc/shadow
/etc/e2guardian/lists/blacklists/drugs/domains
/etc/e2guardian/lists/blacklists/malicious/domains
/etc/rc.common
/etc/uci-defaults/10_migrate-shadow
/lib/preinit/10_indicate_failsafe
/etc/e2guardian/lists/blacklists/sexuality/urls
/init
/etc/e2guardian/lists/blacklists/gambling/expressions
/etc/itus-setup.sh_ORIGINAL
/usr/share/e2guardian/languages/ukenglish/template-itus.html
/lib/preinit/80_mount_root
/etc/rc.local_SLAVE
/etc/rc.local
/lib/functions/uci-defaults.sh
/etc/init.d/led
/lib/functions/network.sh
/etc/init.d/sysctl
/etc/rc.local_BRIDGE
/sbin/hotplug-call
/usr/sbin/csamd
/etc/e2guardian/lists/blacklists/weapons/domains
/etc/e2guardian/lists/blacklists/malicious/expressions
/etc/e2guardian/lists/blacklists/proxies/urls
/etc/e2guardian/lists/blacklists/malware/urls
/etc/rc.button/failsafe
/etc/hotplug.d/iface/15-teql
/usr/sbin/fw_upgrade
/etc/config/dhcp-master
/etc/e2guardian/lists/blacklists/phishing/domains
/www/demo/yday_usage.htm
/etc/e2guardian/lists/blacklists/blasphemy/urls
/etc/services
/www/demo/nmap.html
/lib/functions/system.sh
/bin/board_detect
/etc/uci-defaults/01_network
/lib/upgrade/common.sh
/etc/e2guardian/lists/blacklists/piracy/urls
/etc/config/tinyproxy-gateway
/usr/sbin/ip
/etc/rc.button/power
/etc/sysupgrade.conf
/etc/e2guardian/lists/blacklists/piracy/expressions
/etc/rc.button/rfkill
/lib/upgrade/keep.d/base-files
/etc/init.d/csamd
/lib/preinit/40_run_failsafe_hook
/etc/group
/etc/rc.button/reset
/lib/preinit/30_failsafe_wait
/sbin/firstboot
/www/demo/nmap_sample.html
/usr/sbin/wrtbwmon
/bin/login.sh
/etc/hotplug.d/net/00-sysctl
/etc/rc.local_ROUTER
/etc/hosts
/lib/functions.sh
/etc/e2guardian/lists/blacklists/porn/domains
/etc/iproute2/rt_tables
/etc/uci-defaults/11_migrate-sysctl
/lib/functions/leds.sh
/rom/note
/etc/config/network-bridge
/etc/inittab
/etc/config/network-gateway
/etc/openwrt_release
/sbin/wifi
/etc/profile
/usr/lib/libaivse000.so.5
/etc/itus-setup.sh_BRIDGE
/lib/preinit/99_10_failsafe_login
/etc/init.d/sysfixtime
/etc/TZ
/opt/ctch/etc/csamd.conf
/lib/functions/service.sh
/etc/init.d/umount
/etc/e2guardian/lists/blacklists/racism/domains
/etc/e2guardian/lists/blacklists/gambling/domains
/lib/upgrade/keep.d/base-files-essential
/usr/lib/engines/libocteon.so
/etc/e2guardian/lists/blacklists/spyware/domains
/etc/passwd
/lib/preinit/02_default_set_state
/etc/e2guardian/lists/blacklists/weapons/urls
/etc/e2guardian/lists/blacklists/gambling/urls
/etc/mtab
/etc/e2guardian/lists/blacklists/blasphemy/domains
/etc/init.d/boot
/etc/itus/index.1
/etc/itus/index.2
/lib/uboot-envtools.sh
/etc/e2guardian/lists/blacklists/blasphemy/expressions
/lib/functions/preinit.sh
/lib/upgrade/keep.d/uboot-envtools
/usr/sbin/tc
/bin/ipcalc.sh
/etc/config/tinyproxy-router
/etc/diag.sh
/etc/e2guardian/lists/blacklists/porn/urls
/etc/e2guardian/lists/blacklists/proxies/expressions
/etc/e2guardian/lists/blacklists/spyware/urls
/etc/rc.local_MASTER
/etc/rc.local_GATEWAY
