# usr/lib/opkg/info/snort.list
/usr/lib/snort_dynamicpreprocessor/libsf_sdf_preproc.so.0.0.0
/etc/snort/profiles/suricata/config3_preprocessors.conf
/usr/lib/snort_dynamicpreprocessor/libsf_dce2_preproc.so.0.0.0
/etc/snort/profiles/default/config2_engine.conf
/usr/lib/snort_dynamicpreprocessor/libsf_imap_preproc.so
/usr/lib/snort_dynamicpreprocessor/libsf_dce2_preproc.so.0
/usr/lib/snort_dynamicpreprocessor/libsf_sdf_preproc.so
/etc/snort/profiles/suricata/config5_rules.conf
/etc/snort/preproc_rules/preprocessor.rules
/etc/snort/rules/local.rules
/etc/snort/profiles/snort/config1_advanced.conf
/usr/lib/snort_dynamicpreprocessor/libsf_smtp_preproc.so.0
/usr/lib/snort_dynamicpreprocessor/libsf_sip_preproc.so.0.0.0
/usr/lib/snort_dynamicpreprocessor/libsf_modbus_preproc.so.0.0.0
/etc/snort/gen-msg.map
/usr/lib/snort_dynamicpreprocessor/libsf_ftptelnet_preproc.so
/usr/lib/snort_dynamicpreprocessor/libsf_sdf_preproc.so.0
/etc/snort/profiles/default/config4_other.conf
/etc/snort/profiles/snort/config2_engine.conf
/usr/lib/snort_dynamicpreprocessor/libsf_gtp_preproc.so.0.0.0
/etc/snort/profiles/suricata/config2_engine.conf
/usr/lib/snort_dynamicengine/libsf_engine.so.0
/usr/lib/snort_dynamicpreprocessor/libsf_ssl_preproc.so.0
/usr/lib/snort_dynamicpreprocessor/libsf_dns_preproc.so
/usr/lib/snort_dynamicpreprocessor/libsf_reputation_preproc.so.0.0.0
/etc/snort/preproc_rules/sensitive-data.rules
/usr/lib/snort_dynamicengine/libsf_engine.so
/etc/init.d/snort
/etc/snort/profiles/snort/config5_rules.conf
/etc/snort/rules/snort.rules
/usr/lib/snort_dynamicpreprocessor/libsf_dnp3_preproc.so.0.0.0
/etc/snort/threshold.conf
/etc/snort/rules/L2-sid-msg.map
/etc/snort/attribute_table.dtd
/lib/upgrade/keep.d/snort
/etc/snort/reference.config
/etc/snort/rules/preprocessor.rules
/usr/bin/u2spewfoo
/usr/bin/u2boat
/usr/lib/snort_dynamicpreprocessor/libsf_sip_preproc.so
/usr/lib/snort_dynamicpreprocessor/libsf_ssh_preproc.so
/usr/lib/snort_dynamicpreprocessor/libsf_sip_preproc.so.0
/usr/lib/snort_dynamicpreprocessor/libsf_pop_preproc.so
/etc/snort/unicode.map
/usr/lib/snort_dynamicpreprocessor/libsf_dnp3_preproc.so.0
/etc/snort/preproc_rules/decoder.rules
/etc/snort/profiles/suricata/config1_advanced.conf
/usr/lib/snort_dynamicpreprocessor/libsf_ssh_preproc.so.0.0.0
/usr/lib/snort_dynamicpreprocessor/libsf_ssh_preproc.so.0
/usr/lib/snort_dynamicpreprocessor/libsf_smtp_preproc.so
/etc/snort/profiles/suricata/config4_other.conf
/usr/lib/snort_dynamicpreprocessor/libsf_dns_preproc.so.0.0.0
/usr/lib/snort_dynamicpreprocessor/libsf_gtp_preproc.so.0
/usr/lib/snort_dynamicpreprocessor/libsf_ssl_preproc.so.0.0.0
/etc/snort/profiles/default/config5_rules.conf
/usr/lib/snort_dynamicpreprocessor/libsf_reputation_preproc.so.0
/usr/lib/snort_dynamicpreprocessor/libsf_gtp_preproc.so
/usr/lib/snort_dynamicpreprocessor/libsf_pop_preproc.so.0
/etc/snort/rules/L2.blacklist
/etc/snort/profiles/snort/config3_preprocessors.conf
/usr/lib/snort_dynamicpreprocessor/libsf_reputation_preproc.so
/usr/lib/snort_dynamicpreprocessor/libsf_modbus_preproc.so.0
/usr/lib/snort_dynamicpreprocessor/libsf_pop_preproc.so.0.0.0
/usr/bin/snort
/usr/lib/snort_dynamicengine/libsf_engine.so.0.0.0
/usr/lib/snort_dynamicpreprocessor/libsf_ftptelnet_preproc.so.0.0.0
/etc/snort/exclude_rule.sh
/usr/lib/snort_dynamicpreprocessor/libsf_imap_preproc.so.0
/etc/snort/classification.config
/etc/snort/profiles/snort/config4_other.conf
/usr/lib/snort_dynamicpreprocessor/libsf_ssl_preproc.so
/usr/lib/snort_dynamicpreprocessor/libsf_smtp_preproc.so.0.0.0
/usr/lib/snort_dynamicpreprocessor/libsf_dns_preproc.so.0
/usr/lib/snort_dynamicpreprocessor/libsf_modbus_preproc.so
/usr/lib/snort_dynamicpreprocessor/libsf_ftptelnet_preproc.so.0
/etc/snort/profiles/default/config3_preprocessors.conf
/usr/lib/snort_dynamicpreprocessor/libsf_imap_preproc.so.0.0.0
/etc/snort/rules/decoder.rules
/etc/snort/rules/sensitive-data.rules
/usr/lib/snort_dynamicpreprocessor/libsf_dnp3_preproc.so
/etc/config/snort
/etc/snort/profiles/default/config1_advanced.conf
/usr/lib/snort_dynamicpreprocessor/libsf_dce2_preproc.so
