# usr/lib/opkg/info/strongswan.list
/etc/init.d/ipsec
/usr/lib/ipsec/libhydra.so.0
/usr/lib/ipsec/libstrongswan.so.0
/usr/lib/ipsec/libstrongswan.so.0.0.0
/etc/ipsec.secrets
/etc/strongswan.conf
/usr/lib/ipsec/libhydra.so.0.0.0
/lib/upgrade/keep.d/strongswan
