# usr/lib64/opkg/info/e2guardian.list
/etc/e2guardian/authplugins/ip.conf
/etc/e2guardian/e2guardian.conf
/etc/e2guardian/lists/greysslsitelist
/etc/e2guardian/write_categories.sh
/etc/e2guardian/authplugins/port.conf
/etc/e2guardian/lists/localbannedsearchlist
/etc/e2guardian/lists/blacklists/porn/expressions
/etc/e2guardian/lists/phraselists/games/weighted
/etc/e2guardian/lists/bannedmimetypelist
/etc/e2guardian/lists/contentscanners/exceptionvirusmimetypelist
/usr/share/e2guardian/languages/ukenglish/template.html
/etc/e2guardian/lists/phraselists/violence/weighted_portuguese
/etc/e2guardian/lists/refererexceptionurllist
/etc/e2guardian/lists/phraselists/legaldrugs/weighted
/etc/e2guardian/e2guardianf1.conf
/etc/e2guardian/lists/bannedsitelist
/etc/e2guardian/lists/phraselists/proxies/weighted
/etc/e2guardian/lists/exceptionphraselist
/etc/e2guardian/lists/searchregexplist
/etc/e2guardian/lists/exceptionurllist
/etc/e2guardian/lists/phraselists/pornography/weighted_norwegian
/etc/e2guardian/lists/localexceptionsitelist
/etc/e2guardian/lists/phraselists/pornography/weighted_french
/etc/e2guardian/lists/phraselists/forums/weighted
/etc/e2guardian/lists/phraselists/goodphrases/weighted_general
/etc/e2guardian/lists/phraselists/gore/weighted_portuguese
/etc/e2guardian/lists/phraselists/googlesearches/banned
/etc/e2guardian/lists/exceptionextensionlist
/etc/e2guardian/lists/pics
/etc/e2guardian/lists/phraselists/pornography/weighted_swedish
/etc/e2guardian/lists/phraselists/goodphrases/weighted_general_malay
/etc/e2guardian/lists/phraselists/warezhacking/weighted
/etc/e2guardian/lists/phraselists/pornography/weighted_polish
/etc/e2guardian/lists/phraselists/personals/weighted_portuguese
/etc/e2guardian/lists/phraselists/secretsocieties/weighted
/etc/e2guardian/lists/contentscanners/exceptionvirussitelist
/etc/e2guardian/lists/phraselists/translation/weighted
/etc/e2guardian/lists/authexceptionsitelist
/etc/e2guardian/lists/phraselists/domainsforsale/weighted
/etc/e2guardian/lists/phraselists/pornography/weighted
/etc/e2guardian/lists/localbannedurllist
/etc/e2guardian/lists/exceptioniplist
/etc/e2guardian/lists/phraselists/pornography/weighted_russian
/etc/e2guardian/lists/bannedregexpurllist
/etc/e2guardian/lists/bannedphraselist
/etc/e2guardian/authplugins/ident.conf
/etc/e2guardian/lists/phraselists/goodphrases/weighted_general_danish
/etc/e2guardian/lists/exceptionmimetypelist
/etc/e2guardian/lists/phraselists/illegaldrugs/banned
/etc/e2guardian/lists/phraselists/personals/weighted
/etc/e2guardian/lists/bannediplist
/etc/e2guardian/authplugins/proxy-basic.conf
/etc/e2guardian/downloadmanagers/default.conf
/etc/e2guardian/lists/authexceptionurllist
/etc/e2guardian/lists/greyurllist
/etc/e2guardian/lists/phraselists/badwords/weighted_dutch
/usr/share/e2guardian/languages/ukenglish/messages
/etc/e2guardian/lists/phraselists/pornography/weighted_japanese
/etc/init.d/e2guardian
/etc/e2guardian/lists/phraselists/badwords/weighted_french
/etc/e2guardian/lists/phraselists/news/weighted
/etc/e2guardian/lists/logsitelist
/etc/e2guardian/lists/localbannedsslsitelist
/etc/e2guardian/downloadmanagers/fancy.conf
/etc/e2guardian/lists/phraselists/pornography/banned_portuguese
/etc/e2guardian/lists/localexceptionurllist
/etc/e2guardian/lists/bannedrooms/default
/etc/e2guardian/lists/exceptionsitelist
/etc/e2guardian/lists/logurllist
/etc/e2guardian/lists/phraselists/goodphrases/weighted_general_dutch
/etc/e2guardian/lists/phraselists/pornography/weighted_italian
/etc/e2guardian/lists/downloadmanagers/managedmimetypelist
/etc/e2guardian/lists/phraselists/pornography/weighted_spanish
/etc/e2guardian/lists/headerregexplist
/etc/e2guardian/lists/phraselists/safelabel/banned
/etc/e2guardian/lists/phraselists/goodphrases/weighted_news
/etc/e2guardian/lists/embededrefererurllist
/etc/config/e2guardian
/etc/e2guardian/lists/urlregexplist
/etc/e2guardian/lists/phraselists/gambling/weighted
/etc/e2guardian/lists/phraselists/peer2peer/weighted
/etc/e2guardian/lists/phraselists/illegaldrugs/weighted_portuguese
/etc/e2guardian/lists/phraselists/sport/weighted
/etc/e2guardian/lists/phraselists/conspiracy/weighted
/etc/e2guardian/lists/contentregexplist
/etc/e2guardian/lists/exceptionfileurllist
/etc/e2guardian/installation_or_migration_from_dg.sh
/etc/e2guardian/lists/blacklists/ads/urls
/etc/e2guardian/lists/phraselists/badwords/weighted_spanish
/etc/e2guardian/lists/logregexpurllist
/etc/e2guardian/lists/phraselists/intolerance/banned_portuguese
/etc/e2guardian/lists/phraselists/chat/weighted
/etc/e2guardian/lists/bannedextensionlist
/etc/e2guardian/lists/phraselists/chat/weighted_italian
/etc/e2guardian/authplugins/proxy-digest.conf
/etc/e2guardian/lists/phraselists/badwords/weighted_portuguese
/etc/e2guardian/lists/phraselists/pornography/weighted_malay
/etc/e2guardian/lists/phraselists/gore/weighted
/etc/e2guardian/lists/bannedsearchlist
/etc/e2guardian/lists/authplugins/ipgroups
/etc/e2guardian/lists/phraselists/drugadvocacy/weighted
/etc/e2guardian/contentscanners/clamdscan.conf
/etc/e2guardian/lists/localbannedsitelist
/etc/e2guardian/lists/phraselists/pornography/weighted_danish
/etc/e2guardian/lists/phraselists/intolerance/weighted_portuguese
/etc/e2guardian/lists/exceptionregexpurllist
/etc/e2guardian/lists/phraselists/goodphrases/weighted_general_polish
/etc/e2guardian/lists/refererexceptionsitelist
/etc/e2guardian/lists/bannedsslsitelist
/etc/e2guardian/lists/phraselists/goodphrases/weighted_general_portuguese
/etc/e2guardian/lists/phraselists/pornography/weighted_german
/etc/e2guardian/lists/phraselists/gambling/banned
/etc/e2guardian/lists/phraselists/webmail/weighted
/etc/e2guardian/lists/blacklists/ads/domains
/etc/e2guardian/lists/blacklists/malware/expressions
/etc/e2guardian/lists/addheaderregexplist
/etc/e2guardian/lists/phraselists/music/weighted
/etc/e2guardian/lists/phraselists/pornography/weighted_portuguese
/etc/e2guardian/lists/weightedphraselist
/etc/e2guardian/lists/phraselists/weapons/weighted_portuguese
/etc/e2guardian/lists/phraselists/goodphrases/weighted_general_swedish
/etc/e2guardian/lists/phraselists/goodphrases/exception_email
/usr/sbin/e2guardian
/etc/e2guardian/lists/phraselists/violence/weighted
/etc/e2guardian/lists/phraselists/weapons/weighted
/etc/e2guardian/lists/contentscanners/exceptionvirusurllist
/etc/e2guardian/lists/phraselists/rta/banned
/usr/share/e2guardian/languages/ukenglish/fancydmtemplate.html
/usr/share/e2guardian/blockedflash.swf
/etc/e2guardian/lists/filtergroupslist
/etc/e2guardian/lists/phraselists/nudism/weighted
/etc/e2guardian/lists/localgreysitelist
/etc/e2guardian/lists/phraselists/gambling/weighted_portuguese
/etc/e2guardian/lists/phraselists/intolerance/weighted
/etc/e2guardian/lists/bannedregexpheaderlist
/etc/e2guardian/lists/phraselists/pornography/weighted_russian_utf8
/etc/e2guardian/lists/contentscanners/exceptionvirusextensionlist
/etc/e2guardian/lists/phraselists/idtheft/weighted
/usr/share/e2guardian/transparent1x1.gif
/etc/e2guardian/lists/blacklists/malicious/urls
/etc/e2guardian/lists/embededreferersitelist
/etc/e2guardian/lists/greysitelist
/etc/e2guardian/lists/phraselists/pornography/banned
/etc/e2guardian/lists/localgreyurllist
/etc/e2guardian/lists/bannedurllist
/etc/e2guardian/lists/downloadmanagers/managedextensionlist
/etc/e2guardian/lists/phraselists/goodphrases/exception
/etc/e2guardian/lists/phraselists/illegaldrugs/weighted
/etc/e2guardian/lists/authplugins/portgroups
/etc/e2guardian/lists/urlredirectregexplist
/etc/e2guardian/lists/phraselists/travel/weighted
/etc/e2guardian/lists/phraselists/pornography/weighted_chinese
/etc/e2guardian/lists/phraselists/upstreamfilter/weighted
/etc/e2guardian/lists/exceptionfilesitelist
/etc/e2guardian/lists/localgreysslsitelist
/etc/e2guardian/lists/phraselists/malware/weighted
/etc/e2guardian/lists/phraselists/gambling/banned_portuguese
/etc/e2guardian/lists/phraselists/badwords/weighted_german
/etc/e2guardian/lists/phraselists/pornography/weighted_dutch
