# usr/lib64/opkg/info/ntpclient-noxtal.list
/etc/config/ntpclient
/etc/init.d/ntpclient
/usr/sbin/ntpclient
